# LaunchAgents/com.ejfox.postgres-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ejfox.postgres-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/autossh</string>
        <string>-M</string>
        <string>0</string>
        <string>-N</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <!-- Disable ControlMaster for autossh's ssh child. Without these,
             ssh attaches as master, forks into background, foreground exits
             with status 0, autossh thinks the job is done and exits, launchd
             KeepAlives -> restart loop. Force fresh foreground ssh. -->
        <string>-o</string>
        <string>ControlMaster=no</string>
        <string>-o</string>
        <string>ControlPath=none</string>
        <string>-L</string>
        <string>15432:localhost:5432</string>
        <string>vps</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/postgres-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/postgres-tunnel.log</string>
</dict>
</plist>
