# LaunchAgents/com.ejfox.screenshot-cloudinary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ejfox.screenshot-cloudinary</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ejfox/.dotfiles/bin/screenshot-cloudinary</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Users/ejfox/screenshots</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/ejfox/.local/share/screenshot-cloudinary/watcher.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ejfox/.local/share/screenshot-cloudinary/watcher.log</string>
</dict>
</plist>
